# luxline-monitor.service
[Unit]
Description=LuxLine Website Monitor
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/WebsiteMonitor
Environment="DISCORD_WEBHOOK_URL=your_webhook_url_here"
ExecStart=/home/pi/WebsiteMonitor/.venv/bin/python /home/pi/WebsiteMonitor/main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
